# deploy/roles/libreixp_portal/templates/libreixp_portal.socket.j2
[Unit]
Description=LibreIXP Portal Gunicorn Socket

[Socket]
ListenStream=/run/libreixp_portal.sock

[Install]
WantedBy=sockets.target
