# deploy/roles/teleport/templates/teleport_remove_old_logs.timer.j2
[Unit]
Description=Teleport: Remove Old Logs Every Day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
